# scripts/homeautomation-update.service
[Unit]
Description=Pull latest Home Automation Desktop release from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/homeautomation-scripts/pi-pull-update.sh
